# deploy/inspiration-collector.service
[Unit]
Description=Inspiration Pipeline collector (poll IG collections)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/drew/inspiration-pipeline
ExecStart=/home/drew/inspiration-pipeline/.venv/bin/inspiration collect
