# deploy/inspiration-processor.service
[Unit]
Description=Inspiration Pipeline processor (transcribe + file)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/youruser/inspiration-pipeline
ExecStart=/home/youruser/inspiration-pipeline/.venv/bin/inspiration process
